# root/etc/update-motd.d/98-robot-users
#!/bin/sh
exec sudo -u ros robot users --no-plist
